# systemd/.config/systemd/user/track-pings.service
[Unit]
Description=track pings

[Service]
Type=simple
ExecStart=%h/Dropbox/projects-dropbox/scripts/track_pings

[Install]
WantedBy=default.target
